# production_vibe/avito-ai.service
[Unit]
Description=Avito AI API Server
After=network.target

[Service]
Type=notify
User=zarina
Group=zarina
WorkingDirectory=/home/zarina/Work/BakaiMarket/Avito/production_vibe
Environment="PATH=/home/zarina/Work/BakaiMarket/Avito/production_vibe/venv/bin"
Environment="CUDA_VISIBLE_DEVICES=1"
EnvironmentFile=/home/zarina/Work/BakaiMarket/Avito/production_vibe/.env
ExecStart=/home/zarina/Work/BakaiMarket/Avito/production_vibe/venv/bin/gunicorn -c gunicorn_config.py app_production:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=60
PrivateTmp=true

# Restart policy
Restart=on-failure
RestartSec=10s

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=avito-ai

[Install]
WantedBy=multi-user.target

